# pipeline/systemd/karios-secrets.service
[Unit]
Description=Decrypt Karios secrets into tmpfs (/run/karios) and update trust store
Before=step-ca.service
After=network.target

[Service]
Type=oneshot
Environment=SECRET_STORE=/etc/karios/secrets
Environment=RUNTIME_SECRETS=/run/karios
Environment=STEPCA_USER=root
Environment=STEPCA_GROUP=root
Environment=SOPS_AGE_KEY_FILE=/etc/karios/age/key.txt
ExecStart=/usr/local/bin/karios-secrets-decrypt.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
